# ha-copilot/rootfs/etc/s6-overlay/s6-rc.d/service0_init/up
/etc/s6-overlay/s6-rc.d/service0_init/run
